# docker/elastic/conf/supervisord.minion.conf
[program:salt-minion]
command=/opt/saltstack/salt/salt-minion -l debug
