# docker/elastic/conf/supervisord.master.conf
[program:salt-master]
command=/opt/saltstack/salt/salt-master -l debug
