# docker/elastic/conf/supervisord.loop-jobs.conf
[program:loop-jobs]
command=/bin/sh /usr/bin/loop-jobs.sh
